# backend/crontab
*/5 * * * *  cd /src && /usr/local/bin/python3 tool.delete_demo_users.py
